# scripts/pi-llm-server.service
[Unit]
Description=PI-LLM Unified Server
Documentation=file:///mnt/data0/home/bushuhui_data/msdk/pi-lab/0_ai_agent/pi-llm-server/README.md
After=network.target

[Service]
Type=simple
User=pi-lab
Group=pi-lab
WorkingDirectory=/mnt/data0/home/bushuhui_data/msdk/pi-lab/0_ai_agent/pi-llm-server
Environment="PATH=/home/tiger/anaconda3/envs/vllm/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/home/tiger/anaconda3/envs/vllm/bin/python -m pi_llm_server --config /path/to/config.yaml
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pi-llm-server

# 资源限制
LimitNOFILE=65535
LimitNPROC=4096

# 安全选项
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
